# docker/backup/crontab
# Cultural Sound Lab Backup Schedule
# m h dom mon dow command

# Daily database backup at 2 AM
0 2 * * * /usr/local/bin/backup-database.sh

# Daily storage backup at 3 AM
0 3 * * * /usr/local/bin/backup-storage.sh

# Weekly full system health check at 4 AM on Sundays
0 4 * * 0 /usr/local/bin/backup-health-check.sh

# Clean up old logs daily at 5 AM
0 5 * * * find /backups/logs -name "*.log" -mtime +7 -delete

# Generate backup metrics every hour
0 * * * * /usr/local/bin/backup-metrics.sh